# obsidian-transcriber.service.template
[Unit]
Description=Obsidian Gemini Transcriber Service
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{WORKING_DIR}}
Environment="PATH={{USER_HOME}}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH={{WORKING_DIR}}"
ExecStart={{PYTHON_PATH}} {{WORKING_DIR}}/main.py {{VAULT_PATH}}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target